# packaging/linux/systemd/dasobjectstore-server.service
[Unit]
Description=DASObjectStore standalone Web UI and API
Documentation=https://github.com/sagrudd/DASObjectStore
After=network-online.target dasobjectstored.service
Wants=network-online.target
Requires=dasobjectstored.service

[Service]
Type=simple
Slice=dasobjectstore-control.slice
User=dasobjectstore
Group=dasobjectstore
CPUAccounting=true
MemoryAccounting=true
IOAccounting=true
OOMScoreAdjust=-250
StateDirectory=dasobjectstore
StateDirectoryMode=0750
LogsDirectory=dasobjectstore
LogsDirectoryMode=0750
ExecStart=/usr/bin/dasobjectstore-server --config /opt/dasobjectstore/config.json --generate-missing-tls
Restart=on-failure
RestartSec=5s
UMask=0027
NoNewPrivileges=false
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/run/dasobjectstore /var/lib/dasobjectstore /var/log/dasobjectstore /opt/dasobjectstore

[Install]
WantedBy=multi-user.target
